# app/config/api/crontab
#每分钟处理定时脚本
*/1 * * * * php /home/pplive/web/x-cms/artisan  timing --env=inn>>/tmp/test.log
*/30 * * * * php /home/pplive/web/x-cms/artisan  path --env=inn>>/tmp/path.log